# acs-docker/acs/inserts/supervisord.conf
[supervisord]
nodaemon=true
environment=GENIEACS_CWMP_LOG_FILE="/var/log/genieacs/genieacs-cwmp.log",GENIEACS_CWMP_ACCESS_LOG_FILE="/var/log/genieacs/genieacs-cwmp-access.log",GENIEACS_NBI_LOG_FILE="/var/log/genieacs/genieacs-nbi.log",GENIEACS_NBI_ACCESS_LOG_FILE="/var/log/genieacs/genieacs-nbi-access.log",GENIEACS_FS_LOG_FILE="/var/log/genieacs/genieacs-fs.log",GENIEACS_FS_ACCESS_LOG_FILE="/var/log/genieacs/genieacs-fs-access.log",GENIEACS_UI_LOG_FILE="/var/log/genieacs/genieacs-ui.log",GENIEACS_UI_ACCESS_LOG_FILE="/var/log/genieacs/genieacs-ui-access.log",GENIEACS_DEBUG="true",GENIEACS_DEBUG_FILE="/var/log/genieacs/genieacs-debug.log",GENIEACS_DEBUG_FORMAT="json",GENIEACS_EXT_DIR="/opt/genieacs/ext",GENIEACS_UI_JWT_SECRET="secret",GENIEACS_XMPP_JID="acs@wifi-dev.mig.com.vn",GENIEACS_XMPP_PASSWORD="secret",NODE_OPTIONS="--enable-source-maps",GENIEACS_EXT_TIMEOUT="3000"

[program:genieacs-cwmp]
directory=/opt/genieacs/dist
command=node /opt/genieacs/dist/bin/genieacs-cwmp
stdout_logfile=/var/log/genieacs-cwmp.log
stderr_logfile=/var/log/genieacs-cwmp.log
autorestart=true

[program:genieacs-nbi]
directory=/opt/genieacs/dist
command=node /opt/genieacs/dist/bin/genieacs-nbi
stdout_logfile=/var/log/genieacs-nbi.log
stderr_logfile=/var/log/genieacs-nbi.log
autorestart=true

[program:genieacs-fs]
directory=/opt/genieacs/dist
command=node /opt/genieacs/dist/bin/genieacs-fs
stdout_logfile=/var/log/genieacs-fs.log
stderr_logfile=/var/log/genieacs-fs.log
autorestart=true

[program:genieacs-ui]
directory=/opt/genieacs/dist
command=node /opt/genieacs/dist/bin/genieacs-ui
stdout_logfile=/var/log/genieacs-ui.log
stderr_logfile=/var/log/genieacs-ui.log
autorestart=true

[program:zmq]
directory=/opt/genieacs/dist
command=node /opt/genieacs/ext/zmq.service
stdout_logfile=/var/log/genieacs-ui.log
stderr_logfile=/var/log/genieacs-ui.log
autorestart=true
